# root/etc/uci-defaults/90-prism-migrate-rules
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-only
# Copyright (C) 2026 OctetMill
#
# One-time migration of routing `rule` sections from the pre-redesign single-
# criterion schema (type + value + rs_* + country_*) to the merged-destination
# schema. uci-defaults removes this script once it exits 0.

[ -f /etc/config/prism ] || exit 0

. /lib/functions.sh

changed=0
fakeip_seen=0
provider_seen=""

migrate_rule() {
	local sid rtype uf nm src url label cs n v cc csrc opt
	sid="$1"

	config_get rtype "$sid" type ""
	[ -n "$rtype" ] || return 0
	changed=1

	config_get uf "$sid" use_fakeip "0"
	[ "$uf" = "1" ] && fakeip_seen=1

	case "$rtype" in
	domain|domain_suffix|domain_keyword|domain_regex|ip_cidr|port)
		for v in $(uci -q get "prism.$sid.value" 2>/dev/null); do
			uci add_list "prism.$sid.$rtype=$v"
		done
		;;
	source_ip_cidr)
		# Source matching is dropped — keep the rule disabled and labelled
		# so nothing silently vanishes from the user's config.
		config_get nm "$sid" name ""
		uci set "prism.$sid.enabled=0"
		uci set "prism.$sid.name=${nm} (unsupported - source rule)"
		;;
	remote_ruleset)
		config_get src "$sid" rs_source "sagernet"
		if [ "$src" = "custom" ]; then
			config_get url "$sid" rs_url ""
			label=$(uci -q get "prism.$sid.rs_name" 2>/dev/null | cut -d' ' -f1)
			[ -n "$label" ] || label="$sid"
			if [ -n "$url" ]; then
				cs=$(uci add prism customrs)
				uci set "prism.$cs.label=$label"
				uci set "prism.$cs.url=$url"
				uci add_list "prism.$sid.ruleset=$label"
			fi
		else
			for n in $(uci -q get "prism.$sid.rs_name" 2>/dev/null); do
				uci add_list "prism.$sid.ruleset=$src/$n"
			done
		fi
		;;
	country)
		config_get cc "$sid" country_code ""
		config_get csrc "$sid" country_source "sagernet"
		[ -n "$cc" ] && uci add_list "prism.$sid.country=$cc"
		[ -n "$provider_seen" ] || provider_seen="$csrc"
		;;
	esac

	for opt in type value use_fakeip rs_source rs_delivery rs_name rs_url \
	           rs_format rs_update_interval rs_download_detour \
	           country_code country_source country_delivery; do
		uci -q delete "prism.$sid.$opt"
	done
}

config_load prism
config_foreach migrate_rule rule

[ "$changed" = "1" ] || exit 0

[ "$fakeip_seen" = "1" ] && {
	# A config old enough to predate the merged-destination rule schema may
	# also predate the `dns` section — create it before setting the option.
	uci -q get prism.dns >/dev/null 2>&1 || uci set prism.dns=dns
	uci set prism.dns.fakeip_enabled=1
}
[ -n "$provider_seen" ] && uci set prism.global.country_provider="$provider_seen"

uci commit prism
exit 0
