# root/etc/uci-defaults/92-prism-migrate-subscription-uid
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-only
# Copyright (C) 2026 OctetMill
#
# One-time migration of `subscription` sections to the uid-keyed storage
# scheme. Each subscription gets a stable `uid` UCI option (16 hex chars);
# /etc/prism/nodes/<cfgXXXX>.json is renamed to /etc/prism/nodes/<uid>.json;
# any `node.*.urltest_regex_sources` list that referenced a subscription by
# its cfgXXXX section name is rewritten to reference the new uid.
#
# Background: cfgXXXX anonymous section names are not persisted in
# /etc/config/prism — they are regenerated from a package-wide section counter
# on every UCI parse. Any structural change to the file (delete a manual node,
# drag a section, hand-edit) can rotate them, which silently re-keyed every
# subscription's JSON file out from under its row and broke regex_sources
# lookups. The uid is the stable identifier that pins file storage and the
# regex_sources reference to one specific row for the life of the row.
#
# uci-defaults removes this script once it exits 0.

[ -f /etc/config/prism ] || exit 0

. /lib/functions.sh

NODES_DIR="/etc/prism/nodes"
RENAMES=""

assign_uid() {
	local sid="$1"
	local existing
	config_get existing "$sid" uid
	[ -n "$existing" ] && return 0

	local uid
	uid=$(tr -dc 'a-f0-9' </dev/urandom 2>/dev/null | head -c 16)
	[ -z "$uid" ] && return 1

	uci -q set "prism.$sid.uid=$uid"
	RENAMES="$RENAMES $sid:$uid"
}

config_load prism
config_foreach assign_uid subscription

[ -z "$RENAMES" ] && exit 0

uci -q commit prism

# Rename JSON files in lockstep. The source filename is the current cfgXXXX
# section name; on a router whose anonymous IDs have already rotated past
# their last sync, the file we look for may not exist — that's fine, the
# next sync will create one under the new uid.
for pair in $RENAMES; do
	sid=${pair%%:*}
	uid=${pair#*:}
	src="$NODES_DIR/$sid.json"
	dst="$NODES_DIR/$uid.json"
	if [ -f "$src" ] && [ ! -f "$dst" ]; then
		mv "$src" "$dst"
	fi
done

# Rewrite urltest_regex_sources entries pointing at the old cfgXXXX section
# name to point at the new uid. Values not in the rename map (the `_manual`
# sentinel, an already-uid-shaped value) pass through unchanged.
rewrite_one_node() {
	local nid="$1"
	local current
	current=$(uci -q get "prism.$nid.urltest_regex_sources" 2>/dev/null)
	[ -z "$current" ] && return 0

	local changed=0
	local replacements=""
	local v
	for v in $current; do
		local mapped="$v"
		for pair in $RENAMES; do
			sid=${pair%%:*}
			uid=${pair#*:}
			if [ "$v" = "$sid" ]; then
				mapped="$uid"
				changed=1
				break
			fi
		done
		replacements="$replacements $mapped"
	done

	[ "$changed" = "0" ] && return 0

	uci -q delete "prism.$nid.urltest_regex_sources"
	for v in $replacements; do
		[ -n "$v" ] && uci -q add_list "prism.$nid.urltest_regex_sources=$v"
	done
}

config_load prism
config_foreach rewrite_one_node node
uci -q commit prism

exit 0
