# root/etc/uci-defaults/93-prism-migrate-named-sections
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-only
# Copyright (C) 2026 OctetMill
#
# One-time migration: rename every anonymous UCI section in /etc/config/prism
# to a stable 16-hex-char uid used AS the section name. Eliminates the
# cfgXXXX rotation hazard for every section type (subscription, rule, node,
# condition, customrs, bypass) — not just `subscription` (which decision #14
# / the now-superseded 92 script papered over with a parallel `uid` field).
#
# Affected cross-references, rewritten in lockstep BEFORE the renames:
#   - condition.rule                  (cfg of rule → uid)
#   - node.urltest_regex_sources      (subscription cfg → uid)
#   - /etc/prism/nodes/<sid>.json     (subscription cfg → uid filename)
# A subscription that already ran 92 carries a `uid` UCI field; this script
# adopts that value as its new section name (so files on disk and any
# `urltest_regex_sources` entries already rewritten by 92 stay valid without
# further work) and drops the now-redundant field.
#
# Idempotent: only sections whose current name matches `cfg*` are touched.
# Hand-named sections (user `uci set prism.mysub=subscription`) pass through
# unchanged. Re-runs after a successful run are no-ops.
#
# uci-defaults removes this script once it exits 0.

[ -f /etc/config/prism ] || exit 0

NODES_DIR=/etc/prism/nodes
MAP=$(mktemp) || exit 1
trap 'rm -f "$MAP" "$MAP.rules" "$MAP.subs"' EXIT

gen_uid() {
	tr -dc 'a-f0-9' </dev/urandom 2>/dev/null | head -c 16
}

# Print every section name of the given type. `uci show prism` emits one
# `prism.<sid>=<type>` line per section declaration plus one
# `prism.<sid>.<opt>=<value>` line per option — the section line has
# exactly one dot after `prism`, so n==2 isolates declarations.
section_ids() {
	uci -q show prism | awk -F= -v t="$1" '
		{
			n = split($1, p, ".")
			if (n == 2 && p[1] == "prism" && $2 == t) print p[2]
		}'
}

# Pass 0: build the rename map. Each line: "type old_cfg new_uid".
for type in subscription rule node condition customrs bypass; do
	for sid in $(section_ids "$type"); do
		case "$sid" in
		cfg*) ;;
		*) continue ;;
		esac
		uid=""
		if [ "$type" = "subscription" ]; then
			uid=$(uci -q get "prism.$sid.uid")
		fi
		[ -z "$uid" ] && uid=$(gen_uid)
		[ -z "$uid" ] && continue
		echo "$type $sid $uid" >> "$MAP"
	done
done

[ ! -s "$MAP" ] && exit 0

# Pass 1: rewrite condition.rule. Each condition stores the parent rule's
# old cfg id; remap to the rule's new uid. Done BEFORE any rename so the
# stored value still matches a current section.
awk '$1=="rule" {print $2, $3}' "$MAP" > "$MAP.rules"
if [ -s "$MAP.rules" ]; then
	for csid in $(section_ids condition); do
		rval=$(uci -q get "prism.$csid.rule" 2>/dev/null)
		[ -z "$rval" ] && continue
		new=$(awk -v r="$rval" '$1==r {print $2; exit}' "$MAP.rules")
		[ -n "$new" ] && uci -q set "prism.$csid.rule=$new"
	done
fi

# Pass 2: rewrite node.urltest_regex_sources entries. Each entry is either
# the `_manual` sentinel or a subscription id; remap any subscription id
# that's in the rename map. Already-uid values from 92 won't match the
# cfg-keyed map and pass through unchanged.
awk '$1=="subscription" {print $2, $3}' "$MAP" > "$MAP.subs"
if [ -s "$MAP.subs" ]; then
	for nid in $(section_ids node); do
		current=$(uci -q get "prism.$nid.urltest_regex_sources" 2>/dev/null)
		[ -z "$current" ] && continue
		changed=0
		rebuilt=""
		for v in $current; do
			mapped=$(awk -v r="$v" '$1==r {print $2; exit}' "$MAP.subs")
			if [ -n "$mapped" ]; then
				rebuilt="$rebuilt $mapped"
				changed=1
			else
				rebuilt="$rebuilt $v"
			fi
		done
		[ "$changed" = "0" ] && continue
		uci -q delete "prism.$nid.urltest_regex_sources"
		for v in $rebuilt; do
			[ -n "$v" ] && uci -q add_list "prism.$nid.urltest_regex_sources=$v"
		done
	done
fi

# Pass 3: rename sections. For subscriptions, also rename the on-disk
# node file (skipped on 92-migrated routers — the file is already at
# uid.json) and drop the now-redundant `uid` UCI field.
while read type old new; do
	if [ "$type" = "subscription" ]; then
		src="$NODES_DIR/$old.json"
		dst="$NODES_DIR/$new.json"
		if [ -f "$src" ] && [ ! -f "$dst" ]; then
			mv "$src" "$dst"
		fi
		uci -q delete "prism.$old.uid"
	fi
	uci -q rename "prism.$old=$new"
done < "$MAP"

uci -q commit prism

exit 0
