# root/etc/init.d/prism
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-3.0-only
# Copyright (C) 2026 OctetMill

USE_PROCD=1
START=95
STOP=05

BUILD_CONFIG=/usr/libexec/prism/build-config
FIREWALL=/usr/libexec/prism/firewall.sh
SINGBOX=/usr/bin/sing-box
# Managed-DNS drop-in, generated by build-config from its LOCAL_DOMAINS list.
DNSMASQ_DROPIN_SRC=/var/etc/prism/dnsmasq-managed.conf
# Transient pause marker. Set by the Status page Stop button to suppress
# autostart and the cron watchdog until the user clicks Start (or the
# router reboots — /var/run is tmpfs, so the marker dies with the box and
# an enabled service resumes naturally on next boot).
PAUSE_MARKER=/var/run/prism.paused

# Log a Prism control-plane event to syslog. $1 is the severity
# (info/notice/warn/err); the remaining args form the message.
prism_log() {
	local level
	level="$1"
	shift
	logger -p "daemon.$level" -t prism "$@"
}

install_crons() {
	# Prism drives every periodic task from two cron entries: one hourly
	# maintenance tick (/usr/libexec/prism/hourly — subscription auto-update
	# and the daily rule-set refresh, each self-gated inside that script) and
	# one */5 recovery watchdog. Grep-gated, so this is a no-op once present.
	local changed
	changed=0

	# Sweep the pre-consolidation sync-subscriptions / fetch-rulesets lines
	# (and the pre-jitter fixed `0 * * * *` sync line) left by an older Prism;
	# the hourly tick now subsumes both.
	if grep -qsE 'prism/(sync-subscriptions|fetch-rulesets)' /etc/crontabs/root 2>/dev/null; then
		sed -i -e '/prism\/sync-subscriptions/d' -e '/prism\/fetch-rulesets/d' \
			/etc/crontabs/root
		changed=1
	fi

	# Hourly tick. Random minute 0–59 so a fleet of routers doesn't all hit
	# the subscription / rule-set providers at the same instant. awk (always
	# present on OpenWrt) gives a POSIX-portable PRNG; $RANDOM is a bashism
	# that busybox ash happens to support but POSIX does not.
	grep -qsF "prism/hourly" /etc/crontabs/root 2>/dev/null || {
		local minute
		minute=$(awk 'BEGIN{srand();print int(rand()*60)}')
		echo "$minute * * * * /usr/libexec/prism/hourly" >> /etc/crontabs/root
		changed=1
	}

	# Recovery watchdog every 5 minutes — a downed sing-box must be caught
	# quickly. No jitter: local-only work that contacts no remote server.
	grep -qsF "prism/watchdog" /etc/crontabs/root 2>/dev/null || {
		echo "*/5 * * * * /usr/libexec/prism/watchdog" >> /etc/crontabs/root
		changed=1
	}

	[ "$changed" = 1 ] && /etc/init.d/cron reload 2>/dev/null
}

# Resolve the conf-dir the running dnsmasq instance loads its drop-ins from.
# OpenWrt's dnsmasq.init computes it deterministically, per instance, as
#   config_get dnsmasqconfdir "$cfg" confdir "/tmp/dnsmasq${cfg:+.$cfg}.d"
# Mirror that exactly: derive it from the dnsmasq UCI section name plus an
# explicit `confdir` override if the user set one. The previous version
# parsed the conf-dir= line out of dnsmasq's *generated*
# /var/etc/dnsmasq.conf.<section> instead — which at boot is unreliable: it
# may not exist yet, be half-written by an in-flight dnsmasq restart, or be
# read before its conf-dir= line. Every such failure fell back to
# /tmp/dnsmasq.d, a directory dnsmasq on 25.12 never reads, so the drop-in
# landed where it had no effect and dnsmasq kept answering LAN clients from
# the WAN/ISP (DNS-poisoned) resolvers. UCI is always readable, so deriving
# the path from it works regardless of dnsmasq's runtime state.
dnsmasq_confdir() {
	local section confdir
	section=$(uci -q show dhcp.@dnsmasq[0] 2>/dev/null \
		| sed -n '1s/^dhcp\.\([^.=]*\).*/\1/p')
	[ -n "$section" ] || { echo /tmp/dnsmasq.d; return; }
	confdir=$(uci -q get "dhcp.$section.confdir" 2>/dev/null)
	[ -n "$confdir" ] || confdir="/tmp/dnsmasq.$section.d"
	echo "$confdir"
}

install_dns_dropin() {
	# Managed DNS: hand dnsmasq its forwarding instructions for sing-box's
	# dns-in inbound. The drop-in must land in the conf-dir dnsmasq actually
	# loads (see dnsmasq_confdir). The generated source is deterministic, so
	# an identical copy already in place needs no dnsmasq restart (avoids
	# churn on every reload trigger). build-config writes the source earlier
	# in start_service, so it is present by the time we run.
	[ -f "$DNSMASQ_DROPIN_SRC" ] || return 0

	local confdir dropin
	confdir=$(dnsmasq_confdir)
	dropin="$confdir/prism.conf"

	# Sweep a drop-in left in the wrong directory by an older Prism.
	[ "$confdir" = /tmp/dnsmasq.d ] || rm -f /tmp/dnsmasq.d/prism.conf

	mkdir -p "$confdir"
	cmp -s "$DNSMASQ_DROPIN_SRC" "$dropin" 2>/dev/null && return 0
	cp "$DNSMASQ_DROPIN_SRC" "$dropin"
	# conf-dir is read only at process start, so a reload won't pick this up.
	/etc/init.d/dnsmasq restart >/dev/null 2>&1
	prism_log info "dns: managed DNS drop-in installed in $confdir, dnsmasq routed through sing-box"
}

remove_dns_dropin() {
	local confdir dropin removed
	confdir=$(dnsmasq_confdir)
	removed=0
	for dropin in "$confdir/prism.conf" /tmp/dnsmasq.d/prism.conf; do
		[ -f "$dropin" ] && { rm -f "$dropin"; removed=1; }
	done
	[ "$removed" = 1 ] || return 0
	/etc/init.d/dnsmasq restart >/dev/null 2>&1
	prism_log info "dns: managed DNS drop-in removed, dnsmasq restored to its own resolvers"
}

# Declare the procd instance for sing-box, running it against the config file
# given as $1. Factored out so the config-error fallback in start_service can
# re-declare the existing instance — leaving the running process untouched —
# without duplicating the parameter list.
prism_declare_instance() {
	# Named `singbox` so service_running() in luci.prism and the cron
	# watchdog can query this instance specifically via ubus's `instance`
	# filter, rather than substring-grepping for `"running":true` in the
	# combined service list — which would also match active_watch and
	# defeat the watchdog's whole purpose when clash API is enabled.
	procd_open_instance singbox
	procd_set_param command "$SINGBOX" run -c "$1"
	# procd respawns a crashed sing-box; after a burst of fast crashes it
	# hits its retry limit and gives up. /usr/libexec/prism/watchdog, run
	# from cron, then re-arms the service.
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	# Watch the generated file so an out-of-band edit makes procd restart it.
	procd_set_param file "$1"
	procd_close_instance
}

# Second procd instance: the active-node watcher. Polls sing-box's loopback
# clash API at /proxies every 10s and logs every urltest/selector flip to
# syslog — the diagnostic trail a user needs to decide whether their
# `tolerance` / `interval` are tuned too aggressively. Only declared when
# the clash API itself is enabled (otherwise the daemon has nothing to talk
# to). Restarts in lockstep with sing-box: when the UCI flag flips off, a
# reload re-enters start_service, this function returns without declaring
# the instance, and procd's instance reconciliation tears the daemon down.
prism_declare_active_watch() {
	local clash_api_enabled
	config_get_bool clash_api_enabled global clash_api_enabled 0
	[ "$clash_api_enabled" = "1" ] || return 0
	procd_open_instance active_watch
	procd_set_param command /usr/libexec/prism/active-watch
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	# Hash the script itself, mirroring what we do for sing-box's
	# config file. Without this, procd treats every reload as "same
	# instance, no work to do" — even after the package upgrade
	# rewrites the script on disk. The OLD Lua interpreter then keeps
	# running indefinitely, silently shipping pre-upgrade behaviour.
	# With `file` set, an apk upgrade triggers an automatic restart
	# the next time start_service runs (which happens on every UCI
	# commit via the reload trigger).
	procd_set_param file /usr/libexec/prism/active-watch
	procd_close_instance
}

start_service() {
	local enabled config_path mode cfg_warn chk tmp_cfg config_ok level

	config_load prism
	config_get_bool enabled global enabled 0

	# The UCI `enabled` flag is the persistent master switch: 0 means Prism
	# is installed but dormant — start_service opens no procd instance and
	# the service stays stopped across reboots. Flipped from the UI by the
	# Status page Enable toggle (set_enabled RPC); reload_service routes a
	# disabled service to stop.
	[ "$enabled" = "1" ] || return 1

	# Transient pause: the Status page Stop button writes PAUSE_MARKER and
	# calls /etc/init.d/prism stop, so the service stops without flipping
	# the persistent `enabled` flag. Honour the marker here too — any later
	# reload trigger (UCI commit, WAN-up) re-enters start_service and would
	# otherwise undo the pause. The marker is on tmpfs, so a reboot clears
	# it and an enabled-but-paused service comes back automatically.
	if [ -e "$PAUSE_MARKER" ]; then
		prism_log info "service: paused by user, not starting (clear via Status page Start)"
		return 1
	fi

	# Defer the first start until there is a route out. sing-box's remote
	# rule-set initialization and urltest probes need connectivity, and at
	# cold boot the WAN (DHCP/PPPoE) is not up yet when /etc/rc.d/S95prism
	# runs — starting then leaves sing-box routeless, flooding "no route to
	# internet". With no default route and sing-box not already running,
	# return without an instance: the `interface.*.up wan` trigger re-runs
	# start_service once the WAN comes up (one clean start), and the */5
	# watchdog cron is a backstop. The pgrep guard limits the wait to the
	# genuine first start, so a Save during a WAN outage never stops a
	# healthy instance.
	if ! pgrep sing-box >/dev/null 2>&1; then
		if [ -z "$(ip -4 route show default 2>/dev/null)" ] && \
		   [ -z "$(ip -6 route show default 2>/dev/null)" ]; then
			prism_log info "service: no default route yet, waiting for WAN"
			return 1
		fi
	fi

	config_get config_path global config_path /var/etc/prism/sing-box.json
	config_get mode inbounds mode tproxy

	# Regenerate the sing-box JSON to a side path and validate it there, then
	# promote it over the live config_path only once it is known good. A
	# rejected config therefore never overwrites the last working one — and,
	# crucially, the failure path can still declare a procd instance (see
	# below) rather than returning empty-handed and having procd stop a
	# sing-box that is currently running fine.
	tmp_cfg="${config_path}.next"
	config_ok=1
	if [ -x "$BUILD_CONFIG" ]; then
		# build-config's diagnostics go to stderr; route them into syslog so
		# degraded states (invalid extra.json, dropped DNS rules, WAN-down
		# DNS) stay visible.
		cfg_warn=$("$BUILD_CONFIG" "$tmp_cfg" 2>&1 >/dev/null) || config_ok=0
		if [ -n "$cfg_warn" ]; then
			if [ "$config_ok" = "1" ]; then level=warn; else level=err; fi
			echo "$cfg_warn" | sed 's/^prism build-config: /config: /' \
				| logger -p "daemon.$level" -t prism
		fi
	else
		config_ok=0
		prism_log err "config: build-config is missing, cannot generate sing-box config"
	fi

	# Validate the generated config before handing it to procd. sing-box runs
	# under procd `respawn`, so a config it rejects (e.g. a malformed DNS
	# server address) would otherwise crash-loop forever.
	if [ "$config_ok" = "1" ] && [ -x "$SINGBOX" ]; then
		chk=$("$SINGBOX" check -c "$tmp_cfg" 2>&1) || {
			config_ok=0
			[ -n "$chk" ] && echo "$chk" | logger -p daemon.err -t prism
		}
	fi

	if [ "$config_ok" != "1" ]; then
		rm -f "$tmp_cfg"
		# The freshly generated config is unusable. Returning now without a
		# procd instance would make procd reconcile an empty instance set and
		# STOP a sing-box that is running fine — so a single bad Save would
		# take the LAN down. Instead, when a last working config_path exists,
		# keep procd pointed at it: re-declaring the unchanged instance leaves
		# a running process exactly as it is (same command and file hash, so
		# procd does nothing). The firewall / DNS / cron steps are skipped so
		# the environment stays consistent with the config still in use. Only
		# an initial start with nothing usable on disk genuinely refuses —
		# config_path lives on tmpfs, so at boot it is always absent and a
		# broken config correctly prevents startup.
		if [ ! -f "$config_path" ]; then
			prism_log err "config: generation failed and no working config exists, not starting sing-box"
			return 1
		fi
		prism_log err "config: generated config was rejected, keeping sing-box on the last working config"
		prism_declare_instance "$config_path"
		# Intentionally NOT declaring active_watch here: the last-working
		# config may pre-date the user's clash_api_enabled flip and have
		# no clash listener bound. Polling 127.0.0.1:9090 every 10s
		# against nothing is pure log/CPU noise. The daemon will come
		# back on the next successful Save.
		return 0
	fi

	# Promote the validated config over the live one. The rename is atomic
	# (same tmpfs directory); procd's `file` watch hashes the content, so an
	# unchanged config still yields no restart. Compare before the rename to
	# record which case applies — a Save that produces an identical config is
	# otherwise indistinguishable from one that actually restarted sing-box.
	if [ ! -f "$config_path" ]; then
		prism_log notice "config: writing initial sing-box config ($mode mode)"
	elif cmp -s "$tmp_cfg" "$config_path"; then
		prism_log info "config: sing-box config unchanged, no restart"
	else
		prism_log notice "config: sing-box config changed, restarting sing-box"
	fi

	mv "$tmp_cfg" "$config_path"
	prism_declare_instance "$config_path"
	prism_declare_active_watch

	# Install TPROXY nftables rules + policy routing. No-op for tun mode.
	[ -x "$FIREWALL" ] && "$FIREWALL" start

	# Managed DNS: hand dnsmasq over to sing-box. The generated config has
	# already passed `sing-box check`, so the drop-in can be installed
	# unconditionally; when disabled, clear any stale drop-in.
	local managed_dns
	config_get_bool managed_dns dns managed_dns 1
	if [ "$managed_dns" = "1" ]; then
		install_dns_dropin
	else
		remove_dns_dropin
	fi

	# Install the hourly maintenance tick and the recovery-watchdog cron.
	# Grep-gated, so this is a no-op once the lines are present.
	install_crons

	# Rule-sets are sing-box `remote` rule-sets — sing-box downloads and
	# refreshes them itself, so there is nothing to fetch here.

	return 0
}

stop_service() {
	prism_log info "service: stopping sing-box"

	[ -x "$FIREWALL" ] && "$FIREWALL" stop

	# Return dnsmasq to its own resolvers before sing-box's dns-in goes away.
	remove_dns_dropin

	# Clear the active-node snapshot so the Status page's Groups panel
	# doesn't keep displaying last-known members and latencies after the
	# daemon writing them has been killed by procd. Without this the
	# panel would show stale data with no "service is down" cue.
	rm -f /var/etc/prism/.active-nodes.json

	# Cron lines are left in place. The watchdog respects PAUSE_MARKER and
	# the persistent `enabled` flag, so it self-suppresses for transient
	# pauses and for an explicitly-disabled service. Removing the watchdog
	# cron here defeats its purpose: a manually-stopped service with
	# enabled=1 would have no path back to running until someone calls
	# `start` again. The hourly tick is similarly safe — sync-subscriptions
	# only updates the local cache, and reload_if_changed routes through
	# reload_service which honours both gates. Package removal is the
	# right place to delete the lines, not service stop.
}

reload_service() {
	# Fired by the procd reload trigger (every UCI commit) and by the WAN-up
	# interface trigger declared in service_triggers. Route through start/stop
	# so the `enabled` flag stays authoritative: a Save that turns Prism off
	# tears sing-box down, one that leaves it on re-applies the config.
	#
	# The apply path is a plain `start`, not `stop; start`: start_service is
	# idempotent — install_sync_cron and install_dns_dropin are grep/cmp-gated
	# no-ops when already current, and procd restarts sing-box only when the
	# generated JSON actually changed (the `file` param). So a Save that does
	# not change the config is side-effect-free, without the dnsmasq and
	# sync-cron churn a full stop; start would cause.
	local enabled
	config_load prism
	config_get_bool enabled global enabled 0
	if [ "$enabled" = "1" ]; then
		start
	else
		stop
	fi
}

service_triggers() {
	# React to settings Saves and to the WAN coming up. The interface trigger
	# replaces the old standalone /etc/hotplug.d/iface/40-prism script: procd
	# owns it and registers it only as part of this service, so it cannot
	# fire — and be refused — before the boot-time start, the way a
	# free-standing hotplug script raced S95prism.
	procd_add_reload_trigger prism
	procd_add_interface_trigger "interface.*.up" wan /etc/init.d/prism reload
}
